# config/hosts/koakuma/podman/quadlets/kanidm/kanidm-server.container
[Container]
Image=docker.io/kanidm/server:1.8.5
HostName=idm.koakuma

[Service]
Restart=always

[Install]
WantedBy=default.target